# hidromed.conf
[program:hidromed]
command=/root/hidromed/runserver.sh
ser=www-data
stdout_logfile=/var/log/supervisor/supervisor.stdout.log
stderr_logfile=/var/log/supervisor/supervisor.stderr.log

